# deploy/seoblue0342.service
[Unit]
Description=seoblue0342 - 이후 소설가 네이버 SEO 웹 (Flask/gunicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/seoblue0342
Environment=SEO_DATA_DIR=/opt/seoblue0342/data
# 구글 SERP 노출 체크용 Serper.dev API 키 (발급 절차: deploy/README.md 참고)
# 실제 값을 넣고 주석 해제 후 `systemctl daemon-reload && systemctl restart seoblue0342`
# Environment=SERPER_API_KEY=발급받은_API_키
ExecStart=/opt/seoblue0342/.venv/bin/gunicorn -w 2 -b 127.0.0.1:8842 webapp:app --timeout 120
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
