# deploy/seoblue0342-refresh.timer
[Unit]
Description=매주 월요일 SEO 분석 자동 갱신

[Timer]
OnCalendar=Mon *-*-* 00:05:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
