# deploy/seoblue0342-refresh.timer
[Unit]
Description=매일 새벽 4시(KST) SEO 분석 자동 갱신

[Timer]
# 매일 19:00 UTC = 익일 04:00 KST
OnCalendar=*-*-* 19:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
